# ops/systemd/lumina-research.service
[Unit]
Description=Lumina research & prioritize agent
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
User=ops
WorkingDirectory=/opt/linus-watchtower/ops
ExecStart=/opt/linus-watchtower/ops/run-agent.sh research
TimeoutStartSec=3600
